# deploy/systemd/rustchain-miner.service
[Unit]
Description=RustChain Miner
After=rustchain-node.service
Requires=rustchain-node.service

[Service]
Type=simple
User=rustchain
Group=rustchain
WorkingDirectory=/opt/rustchain
ExecStart=/opt/rustchain/venv/bin/python miners/linux/rustchain_living_museum.py
Restart=always
RestartSec=30
StandardOutput=append:/var/log/rustchain/miner.log
StandardError=append:/var/log/rustchain/miner-error.log
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
