# tools/prometheus/rustchain-exporter.service
[Unit]
Description=RustChain Prometheus Metrics Exporter
After=network.target

[Service]
Type=simple
User=rustchain
WorkingDirectory=/opt/rustchain/tools/prometheus
Environment=NODE_URL=https://rustchain.org
Environment=EXPORTER_PORT=9100
ExecStart=/usr/bin/python3 /opt/rustchain/tools/prometheus/rustchain_exporter.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
